# deploy/boss.service
[Unit]
Description=boss — Home Assistant MQTT to Matter bridge
After=network-online.target avahi-daemon.service
Wants=network-online.target avahi-daemon.service
# Built with the `avahi` feature: boss publishes its Matter mDNS records via
# avahi over D-Bus instead of binding UDP 5353 itself. avahi is a hard
# functional dependency (without mDNS the Matter node is unreachable), so we
# `Wants` it (pull it in + start it) and order `After` it.
# mDNS failure is fatal in boss; disable the start rate-limiter so it keeps
# retrying (e.g. if avahi is briefly unavailable) instead of giving up.
StartLimitIntervalSec=0

[Service]
Type=simple
# Runtime config lives in this env file (broker, state dir, RUST_LOG) as
# BOSS_*/RUST_LOG=value pairs; edit it, not this unit. No leading `-`, so a
# missing config file fails the unit loudly instead of silently using defaults.
# systemd expands the BOSS_* vars into boss's normal flags below (`${VAR}` is a
# single argument), so the binary stays plain flags-only — no env-parsing in it.
EnvironmentFile=/etc/default/boss
ExecStart=/usr/local/bin/boss --mqtt-host ${BOSS_MQTT_HOST} --mqtt-port ${BOSS_MQTT_PORT} --state-dir ${BOSS_STATE_DIR}
Restart=on-failure
RestartSec=5

# Runs as root. The avahi mDNS backend authenticates to the system D-Bus, which
# a transient `DynamicUser` id cannot do; root can. boss binds no privileged
# port directly (avahi owns 5353; boss uses Matter UDP 5540), so this is purely
# about D-Bus access. Acceptable in a dedicated single-purpose container.
# Persisted Matter commissioning state (survives restarts/reboots). systemd
# migrates the previous DynamicUser state from /var/lib/private/boss on switch.
StateDirectory=boss

[Install]
WantedBy=multi-user.target
